# luci-app-openclash-assistant/root/etc/uci-defaults/90_luci-openclash-assistant
#!/bin/sh
[ -f /etc/config/openclash-assistant ] || exit 0
uci -q commit openclash-assistant
exit 0
